# com.printmaxx.auto_outbound_cold_outreach_engine_9569.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.printmaxx.auto_outbound_cold_outreach_engine_9569</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3</string>
        <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/AUTOMATIONS/venture_autonomy.py</string>
        <string>--run</string>
        <string>auto_outbound_cold_outreach_engine_9569</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>WorkingDirectory</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt</string>

    <key>StandardOutPath</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/AUTOMATIONS/logs/auto_outbound_cold_outreach_engine_9569.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/AUTOMATIONS/logs/auto_outbound_cold_outreach_engine_9569_err.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.12/bin:/Users/macbookpro/.local/bin</string>
    </dict>
</dict>
</plist>