# com.claude.schedule.auto_outbound_cold_outreach_engine_9569.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.auto_outbound_cold_outreach_engine_9569</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt" &amp;&amp; claude -p "You are the OUTBOUND autonomy agent for PRINTMAXX venture '\''Cold Outreach Engine'\''. Working directory: /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt  Your job runs every 12h. Each cycle: 1. PROSPECT: Find 20+ new prospects matching the venture'\''s ICP. Use web search, directory scraping, LinkedIn-style lookups. Save to AUTOMATIONS/leads/auto_outbound_cold_outreach_engine_9569/prospects.csv 2. QUALIFY: Score each prospect 1-10 on fit, budget signals, urgency. Move 8+ scores to qualified.csv 3. BUILD ASSET: For each qualified prospect, generate a personalized cold email or DM using the 6-question framework. Save drafts. 4. OUTREACH: Queue the top 10 emails/DMs for sending. Log to outreach_log.csv 5. TRACK: Check previous outreach for replies. Update status.  Rules: - All files stay in /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt - Log everything to AUTOMATIONS/logs/venture_auto_outbound_cold_outreach_engine_9569.log - Use real data, not placeholders - If a step fails, skip to next step, don'\''t stop the cycle  AGENTIC EXECUTION: Do NOT stop after one step. Keep going until the task is COMPLETE or you hit a blocker. - Read state -> Plan -> Execute -> Verify -> If more work needed, CONTINUE - If you find sub-tasks while working, do them now, don'\''t defer - Only stop when: (a) task fully complete, (b) human action required, (c) 15+ tool calls made (safety limit) - Write completion status to AUTOMATIONS/agent/autonomy/auto_outbound_cold_outreach_engine_9569/output/report_20260323.md" --dangerously-skip-permissions --model claude-sonnet-4-6 >> "/Users/macbookpro/.claude/logs/auto_outbound_cold_outreach_engine_9569.log" 2>&amp;1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/macbookpro/.local/bin</string>
    </dict>
    <key>StartInterval</key>
    <integer>43200</integer>
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_outbound_cold_outreach_engine_9569.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_outbound_cold_outreach_engine_9569.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>