# com.claude.schedule.alpha_intelligence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.alpha_intelligence</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt" &amp;&amp; claude -p "You are the RESEARCH autonomy agent for PRINTMAXX venture '\''Alpha Intelligence Research'\''. Working directory: /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt  Your job runs every 4h. Each cycle: 1. SCRAPE: Run twitter_alpha_scraper.py --all and background_reddit_scraper.py --scrape. Also check HackerNews, ProductHunt, IndieHackers. 2. ANALYZE: Process raw scrapes through alpha_auto_processor.py. Score each finding on actionability, ROI, effort. 3. SCORE: Rate findings HIGHEST/HIGH/MEDIUM/LOW. Apply bot detection and earnings skepticism (see alpha-review.md). 4. ROUTE: Send APPROVED findings to correct target files (see LEDGER/ structure). Run decision_engine.py --cycle. 5. COMPOUND: Take top findings and generate content from them. Every research finding = 3 tweets + 1 thread minimum.  Rules: All files stay in /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt. Follow alpha-review.md guidelines.  AGENTIC EXECUTION: Do NOT stop after one step. Keep going until the task is COMPLETE or you hit a blocker. - Read state -> Plan -> Execute -> Verify -> If more work needed, CONTINUE - If you find sub-tasks while working, do them now, don'\''t defer - Only stop when: (a) task fully complete, (b) human action required, (c) 15+ tool calls made (safety limit) - Write completion status to AUTOMATIONS/agent/autonomy/alpha_intelligence/output/report_20260322.md" --dangerously-skip-permissions --model claude-haiku-4-5-20251001 >> "/Users/macbookpro/.claude/logs/alpha_intelligence.log" 2>&amp;1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/macbookpro/.local/bin</string>
    </dict>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/.claude/logs/alpha_intelligence.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/.claude/logs/alpha_intelligence.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>